# rpi_deprecated/rpi.service
[Service]
Type=simple
User=root
ExecStart=/app/rpi
WorkingDirectory=/app
Restart=on-failture
[Install]
WantedBy=multi-user.target
